# {{cookiecutter.app_name}}/etc/systemd/system/{{cookiecutter.app_name}}.socket
[Unit]
Description={{cookiecutter.app_name}} socket

[Socket]
ListenStream=0.0.0.0:8001

[Install]
WantedBy=sockets.target